# docker/supervisor/supervisord.conf
# ===================================================================
# KONFIGURASI SUPERVISOR UNTUK APLIKASI TOKO BERAS
# ===================================================================
#
# Supervisor adalah process manager yang mengelola multiple services
# dalam satu container Docker. File ini mengkonfigurasi:
#
# 1. Nginx - Web server untuk melayani HTTP requests
# 2. PHP-FPM - Application server untuk menjalankan kode PHP
# 3. Laravel Queue Workers - Background job processing
# 4. Laravel Scheduler - Cron job management
#
# Semua services akan dijalankan dan dimonitor secara otomatis
# ===================================================================

# ===================================================================
# KONFIGURASI UTAMA SUPERVISOR
# ===================================================================
[supervisord]
nodaemon=true                                    # Jalankan di foreground (untuk Docker)
user=root                                        # User yang menjalankan supervisor
logfile=/var/log/supervisor/supervisord.log     # File log supervisor
pidfile=/var/run/supervisord.pid                # File PID supervisor

# ===================================================================
# NGINX WEB SERVER
# ===================================================================
# Nginx melayani HTTP requests dan static files
# Priority 10 = start setelah PHP-FPM ready
[program:nginx]
command=nginx -g "daemon off;"                  # Jalankan nginx di foreground
autostart=true                                  # Start otomatis saat supervisor start
autorestart=true                                # Restart otomatis jika crash
priority=10                                     # Priority start (lebih tinggi = start belakangan)
stdout_logfile=/dev/stdout                      # Log stdout ke Docker logs
stdout_logfile_maxbytes=0                       # Unlimited log size
stderr_logfile=/dev/stderr                      # Log stderr ke Docker logs
stderr_logfile_maxbytes=0                       # Unlimited log size

# ===================================================================
# PHP-FPM APPLICATION SERVER
# ===================================================================
# PHP-FPM menjalankan kode PHP Laravel
# Priority 5 = start lebih dulu sebelum Nginx
[program:php-fpm]
command=php-fpm -F                              # Jalankan PHP-FPM di foreground
autostart=true                                  # Start otomatis saat supervisor start
autorestart=true                                # Restart otomatis jika crash
priority=5                                      # Priority start (lebih rendah = start duluan)
stdout_logfile=/dev/stdout                      # Log stdout ke Docker logs
stdout_logfile_maxbytes=0                       # Unlimited log size
stderr_logfile=/dev/stderr                      # Log stderr ke Docker logs
stderr_logfile_maxbytes=0                       # Unlimited log size

# ===================================================================
# LARAVEL QUEUE WORKERS
# ===================================================================
# Queue workers memproses background jobs seperti:
# - Email notifications
# - Image processing
# - Report generation
# - Data import/export
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d  # Nama unik untuk setiap worker process
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
# --sleep=3: tunggu 3 detik jika tidak ada job
# --tries=3: retry job maksimal 3 kali jika gagal
# --max-time=3600: restart worker setiap 1 jam (memory management)
autostart=true                                  # Start otomatis saat supervisor start
autorestart=true                                # Restart otomatis jika crash
stopasgroup=true                                # Stop semua child processes
killasgroup=true                                # Kill semua child processes
user=www-data                                   # Jalankan sebagai user www-data
numprocs=2                                      # Jalankan 2 worker processes parallel
redirect_stderr=true                            # Redirect stderr ke stdout
stdout_logfile=/var/www/html/storage/logs/worker.log  # File log untuk workers
stopwaitsecs=3600                               # Tunggu 1 jam sebelum force kill

# ===================================================================
# LARAVEL TASK SCHEDULER
# ===================================================================
# Scheduler menjalankan cron jobs Laravel seperti:
# - Database cleanup
# - Report generation
# - Data synchronization
# - Automated backups
[program:laravel-schedule]
# Jalankan schedule:run setiap 60 detik (simulasi cron)
command=sh -c "while [ true ]; do (php /var/www/html/artisan schedule:run --verbose --no-interaction &); sleep 60; done"
autostart=true                                  # Start otomatis saat supervisor start
autorestart=true                                # Restart otomatis jika crash
user=www-data                                   # Jalankan sebagai user www-data
redirect_stderr=true                            # Redirect stderr ke stdout
stdout_logfile=/var/www/html/storage/logs/schedule.log  # File log untuk scheduler
